# .docker/cron/crontabs/pollmarkets
* * * * * root curl -H "Accept: application/json" -H "Secret: $(cat /secret)" http://web:8000/poll_markets >> /var/log/cron.log 2>&1
